# scripts/crossx-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.crossx.mvp</string>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>__WORKDIR__/server.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOST</key>
    <string>127.0.0.1</string>
    <key>PORT</key>
    <string>8787</string>
    <key>NODE_ENV</key>
    <string>production</string>
__EXTRA_ENV_VARS__
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__WORKDIR__/data/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__WORKDIR__/data/launchd.err.log</string>
</dict>
</plist>
